# scripts/com.inventoryforge.collector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd job: run the collector every 20 minutes.
  Install:
    cp scripts/com.inventoryforge.collector.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.inventoryforge.collector.plist
  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.inventoryforge.collector.plist
  Tail the log:
    tail -f "<repo>/collector.log"
  NOTE: the absolute paths below assume the repo lives where it is now. If you move
  the repo, update the two <string> paths and reload the job.
  IMPORTANT (macOS): launchd cannot run this if the repo is under Desktop/Documents/
  Downloads/iCloud/OneDrive without Full Disk Access — it fails "Operation not permitted".
  See README "Scheduling" for the fix (move the repo, or grant Full Disk Access to /bin/bash).
-->

<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.inventoryforge.collector</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/michaeljacinto/Library/CloudStorage/OneDrive-Personal/Desktop - onedrive/github/inventoryforge/scripts/run_collector.sh</string>
  </array>
  <!-- Runs a few times a day (9am / 1pm / 6pm) instead of every N minutes, to keep the
       headed-browser popups infrequent. Add/remove <dict> entries to change the times. -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
  </array>
  <key>StandardOutPath</key>
  <string>/Users/michaeljacinto/Library/CloudStorage/OneDrive-Personal/Desktop - onedrive/github/inventoryforge/collector.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/michaeljacinto/Library/CloudStorage/OneDrive-Personal/Desktop - onedrive/github/inventoryforge/collector.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
